# templates/redis.service
[Unit]
Description=redis Server

[Service]
TimeoutStartSec=0
Restart=always
ExecStart=/usr/local/bin/redis
[Install]
WantedBy=multi-user.target